# tests/dockerfiles/Dockerfile.debian-13
FROM debian:13

ENV DEBIAN_FRONTEND=noninteractive

RUN apt-get update && \
    apt-get install -y --no-install-recommends \
        sudo \
        curl \
        wget \
        lsb-release \
        ca-certificates && \
    rm -rf /var/lib/apt/lists/*

# MySQL uses native AIO by default, which fails on Docker's overlayfs storage driver.
RUN mkdir -p /etc/mysql/conf.d && \
    printf '[mysqld]\ninnodb_use_native_aio=0\n' > /etc/mysql/conf.d/docker.cnf

# MySQL 8.4 community server on Debian ships only systemd units (no init.d script).
# This wrapper lets `service mysql start/stop/restart/status` work in a non-systemd
# container by invoking mysqld directly via the helpers bundled with the package.
RUN mkdir -p /etc/init.d && cat > /etc/init.d/mysql << 'EOF'
#!/bin/bash
### BEGIN INIT INFO
# Provides:          mysql
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: MySQL Community Server
# Description:       MySQL Community Server wrapper for non-systemd containers
### END INIT INFO

PIDFILE=/var/run/mysqld/mysqld.pid
MYSQLD=/usr/sbin/mysqld
HELPERS=/usr/share/mysql-8.4/mysql-helpers

do_start() {
    if [ ! -x "$MYSQLD" ]; then
        echo "mysqld not found at $MYSQLD" >&2
        return 1
    fi
    # Source helpers for verify_ready / verify_database / get_running
    . "$HELPERS"
    verify_ready ""
    verify_database ""
    if [ "$(get_running)" = "1" ]; then
        echo "MySQL is already running"
        return 0
    fi
    mkdir -p /var/run/mysqld && chown mysql:mysql /var/run/mysqld
    su -s /bin/bash mysql -c "$MYSQLD --daemonize --pid-file=$PIDFILE" 2>&1
    # Wait for MySQL to accept connections (up to 30s)
    local i=0
    while [ $i -lt 30 ]; do
        mysqladmin ping --silent >/dev/null 2>&1 && return 0
        sleep 1
        i=$((i+1))
    done
    echo "MySQL did not start within 30 seconds" >&2
    return 1
}

do_stop() {
    if [ -f "$PIDFILE" ]; then
        local pid
        pid=$(cat "$PIDFILE" 2>/dev/null) || return 0
        kill "$pid" 2>/dev/null || true
        local i=0
        while [ -d "/proc/$pid" ] && [ $i -lt 20 ]; do
            sleep 1
            i=$((i+1))
        done
        rm -f "$PIDFILE"
    fi
    return 0
}

do_status() {
    if [ ! -x "$MYSQLD" ]; then
        return 3
    fi
    . "$HELPERS"
    if [ "$(get_running)" = "1" ]; then
        return 0
    else
        return 1
    fi
}

case "$1" in
    start)              do_start ;;
    stop)               do_stop ;;
    restart|force-reload) do_stop; sleep 1; do_start ;;
    status)             do_status ;;
    *)  echo "Usage: $0 {start|stop|restart|force-reload|status}"; exit 1 ;;
esac
EOF
RUN chmod +x /etc/init.d/mysql

CMD ["/bin/bash"]
